# vendor/github.com/TheCacophonyProject/modemd/_release/modemd.service
[Unit]
Description=Cacophony Project USB modem controller
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/modemd
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
